# sds/etc/scripts/update.cron
0 0 30 5 * /sds-backend/sds/etc/scripts/update.bash
0 0 30 11 * /sds-backend/sds/etc/scripts/update.bash
